# sands/files/sands.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=63

start() {
	SERVICE_DAEMONIZE=1 service_start /usr/sbin/sands start
}

stop() {
	 ps | grep -E "sands/main|/sands" | grep -v grep  | awk '{print $1}' | xargs kill
}
