# ac-apps/nos-zone/files/etc/hotplug.d/iface/15-nos-zone
#!/bin/sh

. /lib/functions.sh
. /lib/functions/network.sh

case "$ACTION" in
	ifup)
		/etc/init.d/nos-zone restart
	;;
	ifdown)
		:
	;;
esac

exit 0
