# database/files/database.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=61
STOP=10
BIN=/usr/sbin/database
start() {
	SERVICE_DAEMONIZE=1 service_start $BIN start 
	echo database done > /tmp/debug
}

stop() {
	ps | grep -E "database/main|/database" | grep -v grep  | awk '{print $1}' | xargs kill 
}
