# proxybase/files/proxybase.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=95

USE_PROCD=1

start_service() {
	procd_open_instance
	procd_set_param command /usr/sbin/proxybase start
	procd_set_param respawn
	procd_close_instance
}

stop_service() {
	ps | grep -E "proxybase/main.lua" | grep -v grep  | awk '{print $1}' | xargs kill
}