# base-config/files/base-init.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2006-2011 OpenWrt.org

START=09

boot() {
	test -f /etc/base_reset && {
		echo -ne 'd\n4\nw\n' | fdisk /dev/sda
		echo -ne 'd\n3\nw\n' | fdisk /dev/sda
		echo -ne 'n\np\n3\n1835008\n3932160\nn\np\n3932672\n\nw\n' | fdisk /dev/sda
		partx -a /dev/sda
		echo -ne 'n\n' | mkfs.ext4 /dev/sda3
		echo -ne 'n\n' | mkfs.ext4 /dev/sda4
		rm -f /etc/base_reset
	}
}
