# ac-apps/nos-tbqd/files/etc/hotplug.d/iface/15-nos-tbqd
#!/bin/sh

. /lib/functions.sh
. /lib/functions/network.sh

case "$ACTION" in
	ifup)
		/etc/init.d/nos-tbqd restart
	;;
	ifdown)
		:
	;;
esac

exit 0
